# mosquitto/services/mqtt/run
#!/bin/bash
exec /sbin/setuser abc /usr/local/sbin/mosquitto -c /config/mosquitto.conf 